# services/llm-serve-stop.timer
# ~/.config/systemd/user/llm-serve-stop.timer
[Unit]
Description=Stop LLM server each evening

[Timer]
OnCalendar=*-*-* 19:00:00 Asia/Taipei
Persistent=true

[Install]
WantedBy=timers.target
