# services/llm-serve-stop.service
# ~/.config/systemd/user/llm-serve-stop.service
[Unit]
Description=Stop LLM server

[Service]
Type=oneshot
ExecStart=systemctl --user stop llm-serve-start.service
